# crontab
0 23 * * * python /usr/src/app/updateFantasyWeekStats.py
0 23 * * * python /usr/src/app/test.py
# empty line required
